# root/etc/services.d/privoxy/run
#!/usr/bin/with-contenv bash
exec privoxy --no-daemon ${PRIVOXY_CONFIG:-/etc/privoxy/config}
